# Understand_init_managers/Create_a_recipe_holding_a_SystemD_service/dirs/meta-tutorial/recipes-init-managers/systemd-example/files/mysysd.service
[Unit]
Description= A continuous systemd service to demonstrate how init managers work

[Service]
ExecStart=/bin/bash /usr/bin/mydaemon.sh

[Install]
WantedBy=multi-user.target

